# .docker/rabbitmq.sv.conf
[program:rabbitmq]
command=/bin/bash -c "runuser -l rabbitmq -c '/usr/sbin/rabbitmq-server'"
redirec_stderr=true
startsecs=5
